# root/etc/s6-overlay/s6-rc.d/svc-queue-worker/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "*** Starting Async Action Queue ***"

exec /usr/bin/php /app/www/artisan queue:work --sleep=3 --tries=1 --max-time=3600
