# root/etc/services.d/memcached/run
#!/usr/bin/with-contenv bash
exec memcache memcached
